# docs/ts-cftunnel-automated/config/supervisord.conf
# config/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:tailscaled]
command=/usr/local/bin/tailscaled -state=/var/lib/tailscale/tailscaled.state -socket=/var/run/tailscale/tailscaled.sock -tun=tailscale0
autostart=%(ENV_ENABLE_TAILSCALE)s
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:tailscale-up]
command=/usr/local/bin/tailscale-up.sh
autostart=%(ENV_ENABLE_TAILSCALE)s
autorestart=unexpected
startsecs=0
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:cloudflared]
command=/usr/local/bin/cloudflared-up.sh
# Only start if enabled AND either token or ID is provided
autostart=%(ENV_ENABLE_CLOUDFLARE)s
autorestart=%(ENV_ENABLE_CLOUDFLARE)s
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0